# stash/rootfs/etc/services.d/stash/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree when Stash stops
# ==============================================================================

/run/s6/basedir/bin/halt
